# systemd/user/rclipboard.service
[Unit]
Description=rclipboard service (user)
After=default.target

[Service]
Type=simple
EnvironmentFile=%h/.config/rclipboard/env
# Set this to your repo path or install path
# Override via: systemctl --user edit rclipboard.service
WorkingDirectory=%h/rclipboard-py
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host ${RCLIPBOARD_BIND_ADDR} --port ${RCLIPBOARD_BIND_PORT}
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target

